# src/aipass/skills/lib/telegram/telegram-bot@.service
# Systemd user service template for AIPass Telegram bots.
#
# Install:
#   cp telegram-bot@.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#
# Usage:
#   systemctl --user start  telegram-bot@<bot_id>
#   systemctl --user stop   telegram-bot@<bot_id>
#   systemctl --user enable telegram-bot@<bot_id>   # auto-start on login
#   systemctl --user status 'telegram-bot@*'        # all bots
#
# The %i specifier is replaced by the bot_id (e.g. "base", "dev_central").

[Unit]
Description=AIPass Telegram Bot — %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/Projects/AIPass/.venv/bin/python3 -m aipass.skills.lib.telegram.apps.handlers.base_bot --bot-id %i
WorkingDirectory=%h/Projects/AIPass
Environment=AIPASS_BOT_ID=%i
Environment=AIPASS_SESSION_TYPE=telegram
Restart=on-failure
RestartSec=10
StandardOutput=append:%h/Projects/AIPass/system_logs/telegram-bot-%i.log
StandardError=append:%h/Projects/AIPass/system_logs/telegram-bot-%i.log

[Install]
WantedBy=default.target
