# src/aipass/prax/prax-monitor.service
# Systemd user service for the prax monitor with Telegram relay.
#
# Install:
#   cp prax-monitor.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#
# Usage:
#   systemctl --user start  prax-monitor
#   systemctl --user enable prax-monitor   # auto-start on login
#   systemctl --user status prax-monitor

[Unit]
Description=AIPass Prax Monitor — Telegram relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# NOTE: the module __main__ takes a single positional (the subcommand); 'run' →
# _run_monitor([]) = all branches. The relay is enabled by the env var below
# (the __main__ argparse rejects extra tokens like 'all --relay').
ExecStart=%h/Projects/AIPass/.venv/bin/python3 -m aipass.prax.apps.modules.monitor run
WorkingDirectory=%h/Projects/AIPass
Environment=AIPASS_PRAX_MONITOR_RELAY=1
Restart=always
RestartSec=5
# IMPORTANT: log OUTSIDE system_logs/ — the monitor tails system_logs/*.log and
# @trigger watches it too. Writing the monitor's own output there creates a
# feedback loop (monitor output → re-tailed/recorded → reported → more output).
StandardOutput=append:%h/.aipass/prax-monitor.service.log
StandardError=append:%h/.aipass/prax-monitor.service.log

[Install]
WantedBy=default.target
