# src/aipass/trigger/templates/trigger-log-watcher.service.template
[Unit]
Description=AIPass Trigger Log Watcher
After=default.target

[Service]
Type=simple
WorkingDirectory={{AIPASS_HOME}}
ExecStart={{AIPASS_HOME}}/.venv/bin/python3 -m aipass.trigger.apps.log_watcher_service
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
